# DEBIAN/etc/systemd/pio2-cient.service
[Unit]
Description=Piot2 Client Service

[Service]
ExecStart=/opt/piot2/piot2-ctrl.sh client-hook
User=__pio2

[Install]
WantedBy=multi-user.target
