# debian/prerm
#!/bin/sh
set -e

systemctl disable piot2-client.service
systemctl disable piot2-client.timer
systemctl disable piot2-server.service
systemctl disable piot2-server.timer
